# plymouth/save-brightness.conf
start on runlevel [06]
task

script
    . /etc/default/chromebook-brightness-hack
    mkdir -p `dirname $LEVEL_SAVE_FILE`
    cat $SYS_BRIGHTNESS_FILE > $LEVEL_SAVE_FILE
end script
